# deploy/supervisor/elevasi-worker.conf
; Elevasi — Laravel queue worker (Redis)
; Copy to: /etc/supervisor/conf.d/elevasi-worker.conf
; Handles inquiry emails + Spatie media conversions.

[program:elevasi-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /home/uXXXXX/domains/elevasidesignbuild.com/public_html/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
directory=/home/uXXXXX/domains/elevasidesignbuild.com/public_html
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=uXXXXX
numprocs=1
redirect_stderr=true
stdout_logfile=/home/uXXXXX/domains/elevasidesignbuild.com/public_html/storage/logs/worker.log
stopwaitsecs=3600
