# deploy/supervisor/elevasi-ssr.conf
; Elevasi — Inertia SSR (Node)
; Copy to: /etc/supervisor/conf.d/elevasi-ssr.conf
; Then: sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start elevasi-ssr:*
;
; Adjust user= and directory= to your VPS paths.

[program:elevasi-ssr]
process_name=%(program_name)s
command=php /home/uXXXXX/domains/elevasidesignbuild.com/public_html/artisan inertia:start-ssr
directory=/home/uXXXXX/domains/elevasidesignbuild.com/public_html
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=uXXXXX
numprocs=1
redirect_stderr=true
stdout_logfile=/home/uXXXXX/domains/elevasidesignbuild.com/public_html/storage/logs/ssr.log
stopwaitsecs=10
